# local.crontab
# local.crontab
# -------------
# format:
# mm hh dd MM yyyy <job>
* * * * * /bin/ls -lash # every minute 
47 * * * * /bin/sleep 2 && echo 'complete' # every hour
# and so on, so on ...
